# cs2-panel.service
[Unit]
Description=CS2 Web Control Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cs2-panel
ExecStart=/usr/bin/node /opt/cs2-panel/server.js
Restart=on-failure
RestartSec=5
# ---- EDIT THESE ----
Environment=PANEL_USER=admin
Environment=PANEL_PASS=ChangeThisPanelPassword
Environment=PANEL_PORT=8080
Environment=SESSION_SECRET=some-long-random-string-here
Environment=VARS_FILE=/home/steam/cs2_server/server-vars.conf
Environment=CS2_SERVICE=cs2
# ---- "Update Server" button (SteamCMD app_update 730) ----
# The panel runs as root, so it can `su - steam` to run SteamCMD as the game user.
Environment=STEAM_USER=steam
Environment=STEAMCMD=/home/steam/steamcmd/steamcmd.sh
Environment=CS2_DIR=/home/steam/cs2_server
# Or override the whole command:
# Environment=UPDATE_CMD=su - steam -c "/home/steam/steamcmd/steamcmd.sh +force_install_dir /home/steam/cs2_server +login anonymous +app_update 730 validate +quit"

[Install]
WantedBy=multi-user.target
